# fc-webhook.service
[Unit]
Description=FC-Advisor Webhook Server
After=network.target

[Service]
Type=simple
User=tely
WorkingDirectory=/home/tely/fc-advisor
Environment=WEBHOOK_SECRET=-jahc0WAj_FOC12U8wuhAME_rJPDLsC9bT9oz7yigVo
ExecStart=/usr/bin/python3 /home/tely/fc-advisor/webhook_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
